# deploy/debt-risk-radar.service
[Unit]
Description=Debt Risk Radar Streamlit service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debt-radar
Group=debt-radar
WorkingDirectory=/opt/debt-risk-radar
EnvironmentFile=/etc/debt-risk-radar.env
Environment=HOME=/var/lib/debt-risk-radar
ExecStart=/opt/debt-risk-radar/.venv/bin/streamlit run /opt/debt-risk-radar/app.py --server.address 127.0.0.1 --server.port 8501 --server.headless true --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/debt-risk-radar
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
